# nvidia-fan-curve-failure.service
[Unit]
Description=Log concise failure reason for nvidia-fan-curve.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nvidia-fan-curve-onfailure.sh
